# googel.service
[Unit]
Description=Googel
Requires=network.target
After=network.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/home/home/Desktop/googel/staticserver --https=false --mode=spa --port=9987
ExecStop=pkill -x 'staticserver'
WorkingDirectory=/home/home/Desktop/googel/
StandardOutput=file:/home/home/Desktop/googel/stdout-file

[Install]
WantedBy=multi-user.target
