# googel.service
[Unit]
Description=Googel
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=on-failure
User=root
ExecStart=/home/home/Desktop/googel/staticserver --mode=spa --port=9987
WorkingDirectory=/home/home/Desktop/googel/
StandardOutput=file:/home/home/Desktop/googel/stdout-file


[Install]
WantedBy=multi-user.target
